# supervisord.conf
[program:node]
command = ./bin/chanjs
directory = /home/dotcloud/current